# doc/runit/shelper/run
#!/bin/sh

exec /u/apps/shelper/bin/shelper-agent
