# doc/runit/shelper/run
#!/bin/sh

exec /usr/bin/shelper-agent /etc/shelper/shelper.yml
